# deploy/supervisor.conf
# Supervisor 进程管理配置 — 多巴胺行情终端
# 适用: Ubuntu/Debian/CentOS
# 放置: /etc/supervisor/conf.d/dopamine.conf
# 重载: supervisorctl reread && supervisorctl update

[program:dopamine_market]
command=/opt/dopamine-market/venv/bin/gunicorn -c deploy/gunicorn.conf.py app:app
directory=/opt/dopamine-market
user=www-data
autostart=true
autorestart=true
startretries=3
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/opt/dopamine-market/logs/supervisor_stdout.log
stderr_logfile=/opt/dopamine-market/logs/supervisor_stderr.log

# 环境变量
environment=FLASK_ENV="production",DB_TYPE="sqlite"
